# user.scripts/sd_import/99-sd-card-import.rules
# Unraid SD Card Auto-Import Rule
# Triggers import script when SD card reader detects media
#
# SETUP: Set SD_CARD_READER in .env file, then run install-sd-import
# Find your serial with: udevadm info --query=all --name=/dev/sdX | grep ID_SERIAL
#
# Note: "change" action catches media insertion when reader stays connected

# Cleanup mounts first on any change (in case card was swapped)
ACTION=="add|change", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", \
  ENV{ID_SERIAL}=="YOUR_SD_READER_SERIAL", \
  RUN+="/bin/bash -c 'for mp in /mnt/sd-import/*; do [ -d \"$mp\" ] && umount \"$mp\" 2>/dev/null && rmdir \"$mp\"; done'"

# Then trigger import if media is present (has partitions)
ACTION=="add|change", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", \
  ENV{ID_SERIAL}=="YOUR_SD_READER_SERIAL", \
  RUN+="/bin/sh -c 'echo \"/usr/local/bin/sd-card-import %k\" | /usr/bin/at -M now'"

ACTION=="remove", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", \
  ENV{ID_SERIAL}=="YOUR_SD_READER_SERIAL", \
  RUN+="/bin/bash -c 'for mp in /mnt/sd-import/*; do [ -d \"$mp\" ] && umount \"$mp\" 2>/dev/null && rmdir \"$mp\"; done; rmdir /mnt/sd-import 2>/dev/null || true'"
